# packages/sources/tmjos-dock/files/autostart/plank.desktop
[Desktop Entry]
Type=Application
Name=Plank
Comment=Stupidly simple dock
Exec=plank
Icon=plank
Categories=GNOME;GTK;Utility;
StartupNotify=false
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Phase=Applications
