# packages/sources/tmjos-dock/files/autostart/tmjos-dock-sync.desktop
[Desktop Entry]
Type=Application
Name=TMJOs Dock Sync
Comment=Sincroniza GNOME favorite-apps com Plank dock-items
Exec=/usr/bin/tmjos-dock-sync
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Phase=Applications
